# deploy/convorocp-terminal.service
[Unit]
Description=ConvoroCP web terminal (ttyd) — serves an interactive PTY on loopback; nginx proxies /tty/ to it, gated to operators via auth_request
After=network.target

[Service]
# Bound to 127.0.0.1 only — never expose ttyd directly. -b /tty matches the
# nginx proxy base path so ttyd's assets resolve under /tty/. -W = writable
# (accept keyboard input). Runs as root → operator gets a root shell on the node.
ExecStart=/usr/bin/ttyd -b /tty -i 127.0.0.1 -p 13715 -W bash
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
